# deploy/uzt-strategy-api.service
[Unit]
Description=Uzt Strategy API
After=network.target

[Service]
WorkingDirectory=/srv/uzt-strategy/backend
ExecStart=/usr/bin/node src/server.js
Restart=always
EnvironmentFile=/srv/uzt-strategy/backend/.env

[Install]
WantedBy=multi-user.target
